# squashfs-root/etc/inittab.default
::sysinit:/etc/init.d/rcS
#::askfirst:-/bin/ash
